# puppet/files/pa-lbaas-agent.upstart
description     "PortAuthority LBaaS Agent"
start on docker and etcd
stop on runlevel [06]
respawn
setuid root
setgid root
console log
exec /usr/local/bin/pa-lbaas-agent
